# systemd/user/ndb_ygg_guard.service
[Unit]
Description=NomadLink-DataBridge guard (ensures Yggdrasil overlay presence, prints status)
After=default.target

[Service]
Type=simple
ExecStart=%h/NomadLink-MIPA/scripts/ndb_status.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
